# CommonCrawl/crontab
*/5 * * * * run CommonCrawl cc-cron.often
