# KB-NL.Fryslan/crontab

10 0 * * * run KB-NL.Fryslan kb-cron.daily
